# systemd/signalmasterpro.service
[Unit]
Description=SignalMasterPro FastAPI app
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/SignalMasterPro
ExecStart=/opt/SignalMasterPro/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
